# systemd/pirelayserver.service
[Unit]
Description=PiRelayServer

[Service]
Type=simple
ExecStart=/opt/pirelayserver/pirelayserver

[Install]
WantedBy=multi-user.target